# deploy/systemd/messagechain-validator.service
[Unit]
Description=MessageChain validator node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=messagechain
Group=messagechain
WorkingDirectory=/opt/messagechain
ExecStart=/usr/bin/python3 /opt/messagechain/server.py \
    --data-dir /var/lib/messagechain \
    --rpc-bind 0.0.0.0 \
    --keyfile /etc/messagechain/keyfile
Environment=MESSAGECHAIN_RPC_AUTH_ENABLED=false
Environment=MESSAGECHAIN_REQUIRE_CHECKPOINTS=false
# TESTING: 30s block time for fast live-network iteration.  Remove this
# line (or flip to 600s) for production-cadence operation.
Environment=MESSAGECHAIN_BLOCK_TIME_TARGET=30
# TESTING: tree_height=16 matches the founder-launch entity baked into
# genesis.  Production should bump launch script + this override to 20
# for longer hot-key runtime (~2y vs ~450d at 600s block cadence).
Environment=MESSAGECHAIN_MERKLE_TREE_HEIGHT=16
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/messagechain
ReadOnlyPaths=/etc/messagechain /opt/messagechain

[Install]
WantedBy=multi-user.target
